# config/site_001.conf
[program:site_001]
command=/var/www/venv/bin/gunicorn site_001.wsgi:application -c /var/www/site_001/config/gunicorn.conf.py
directory=/home/django/site_001
user=django
autorestart=true
redirect_stderr=true
stdout_logfile=/home/django/site_001/logs/debug.log